# docker-limiter-desktop/etc/systemd/system/docker-limiter-desktop.path
[Path]
PathModified=/etc/default/docker-limiter-desktop

[Install]
WantedBy=multi-user.target